# tools/vendor/bundle/ruby/2.7.0/gems/ds9-1.4.1/ext/ds9/tmp/x86_64-linux-gnu/ports/nghttp2/v1.34.0/nghttp2-1.34.0/contrib/nghttpx-upstart.conf
# vim: ft=upstart:

description "HTTP/2 reverse proxy"

start on runlevel [2]
stop on runlevel [016]

exec /home/isucon/webapp/tools/vendor/bundle/ruby/2.7.0/gems/ds9-1.4.1/ext/ds9/ports/x86_64-linux-gnu/nghttp2/v1.34.0/bin/nghttpx
